# stoker.conf
start on (started network-interface
          or started network-manager
          or started networking)

stop on (stopping network-interface
         or stopping network-manager
         or stopping networking)

respawn

# this assumes that there is a stoker user defined in the system and it is running this microservice
setuid stoker
chdir /home/stoker/stoker

script
# this script runs in /bin/sh by default
# respawn as bash so we can source in rbenv
exec /bin/bash <<EOT
  # uncomment to use syslog for logging
  # exec &> /dev/kmsg

  export HOME=/home/stoker
  export PATH="/home/stoker/.rbenv/bin:/home/stoker/.rbenv/shims:$PATH"

  # Pick your poison :) Or none if you're using a system wide installed Ruby.
  # rbenv
  source $HOME/.bashrc
  # OR
  # source /home/apps/.profile
  # OR system:
  # source /etc/profile.d/rbenv.sh
  #
  # rvm
  # source /home/apps/.rvm/scripts/rvm

  exec bundle exec puma -C ./config.rb
EOT
end script

